# apps/gjiwang/build/src/supervisor/conf.d/queue.conf
[program:queue]
process_name=%(program_name)s_%(process_num)02d
command=php /var/web/project/app/artisan queue:work --sleep=3 --tries=3
autostart=true
autorestart=true
user=root
numprocs=8
redirect_stderr=true
stdout_logfile=/var/web/project/app/storage/logs/queue.log
stopwaitsecs=3600